# process/tmp/crontab
# Every minute
* * * * * /usr/local/bin/syncdata.sh all > /dev/stdout
# An empty line is required at the end of this file for a valid cron file.